# infra/systemd/nvidia-exporter-p100.service
[Unit]
Description=NVIDIA GPU Exporter (P100 GPU metrics)
After=network.target

[Service]
# --query-field-names restricts to safe fields; newer nvidia-smi on P100 exposes
# fields with spaces/brackets in names that are invalid Prometheus metric names.
ExecStart=%h/.local/bin/nvidia_gpu_exporter \
  --web.listen-address=:9835 \
  --query-field-names=name,driver_version,gpu_uuid,memory.total,memory.used,memory.free,utilization.gpu,utilization.memory,temperature.gpu,power.draw,clocks.sm,clocks.mem
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
