# infra/hydra-head/hydra-head.service
[Unit]
Description=Hydra Head - Node Agent for llama-server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/hydra
ExecStart=/opt/hydra/bin/hydra-head \
  -global /opt/hydra/config/global.yaml \
  -node /opt/hydra/config/node-p100.yaml \
  -api-port 9700
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hydra-head

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=default.target
